# egroupware.cron
MAILTO=root

*/5 * * * *	http /usr/bin/php /usr/share/egroupware/phpgwapi/cron/asyncwrapper.php
